# docs/examples/sec-ui.service.example
[Unit]
Description=SEC Snippet Control Local UI Daemon
After=network.target

[Service]
Type=simple
# NOTE: Adapt WorkingDirectory to your local repository path, e.g. %h/repos/snippet-engine-control
WorkingDirectory=%h/repos/snippet-engine-control
# Environment="ESPANSO_DIR=%h/.config/espanso"
ExecStart=/usr/bin/npm run ui
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
